# .config/systemd/user/caffeine.service
[Unit]
Description=speed

[Service]
ExecStart=/usr/bin/caffeine
