# .config/systemd/user/fcitx5.service
[Unit]
Description=nihongo syabarenai

[Service]
ExecStart=/usr/bin/fcitx5
