# .config/systemd/user/lock.service
[Unit]
Description=no

[Service]
ExecStart=/home/lily/.scripts/lock.sh
