# .config/systemd/user/mconnect.service
[Unit]
Description=u got games on ur phone?

[Service]
ExecStart=/usr/local/bin/mconnect
