# .config/systemd/user/redshift.service
[Unit]
Description=eye strain

[Service]
ExecStart=/usr/bin/redshift-gtk
