# .config/systemd/user/stalonetray.service
[Unit]
Description=more bullshit

[Service]
ExecStart=/usr/bin/stalonetray
